# templates/Library/LaunchAgents/com.user.backup.icloud.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.backup.icloud</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/rclone</string>
        <string>copy</string>
        <string>--include-from</string>
        <string>{{ansible_env.HOME}}/.config/backup/icloud.conf</string>
        <string>home:</string>
        <string>icloud-crypt:backup/{{lookup('pipe', 'hostname')}}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/dev/null</string>
</dict>
</plist>

