# templates/Library/LaunchAgents/com.user.backup.wasabi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.backup.wasabi</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/rclone</string>
        <string>copy</string>
        <string>--include-from</string>
        <string>{{ansible_env.HOME}}/.config/backup/wasabi.conf</string>
        <string>home:</string>
        <string>wasabi-backup:{{lookup('pipe', 'hostname')}}</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer> <!-- Hour of the day (0 = midnight) -->
        <key>Minute</key>
        <integer>0</integer> <!-- Minute of the hour -->
    </dict>

    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/dev/null</string>
</dict>
</plist>

